# docker2/supervisor/scron.conf
[program:cron]
command = /usr/sbin/cron -f -L 15
autostart=true
autorestart=true
priority=30
