# docker/supervisor/srest.conf
[program:rest]
command = rest
directory=/src
stdout_logfile=/var/log/supervisor/rest-stout
stderr_logfile=/var/log/supervisor/rest-stderr
autostart=false
autorestart=true
priority=30
startsecs=60
