# docker/supervisor/srsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=false
autorestart=true
redirect_stderr=true
priority=10
